# telegram-bridge.service
[Unit]
Description=Telegram Bridge API (Universal File Storage)
After=network.target
Wants=network-online.target

[Service]
Type=notify
User=nhprince
Group=nhprince
WorkingDirectory=/home/nhprince/telegram-bridge
RuntimeDirectory=telegram-bridge
Environment="PATH=/home/nhprince/telegram-bridge/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStartPre=/usr/bin/mkdir -p /var/log/telegram-bridge /tmp/telegram-bridge-uploads
ExecStartPre=/usr/bin/chown nhprince:nhprince /var/log/telegram-bridge /tmp/telegram-bridge-uploads
ExecStart=/home/nhprince/telegram-bridge/venv/bin/gunicorn main:app -c gunicorn.conf.py
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
KillMode=mixed
TimeoutStopSec=30
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening (minimal — cloudflare tunnel + nginx handles DDoS/SSL)
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
